# manifests/maverick-modules/maverick_mavlink/files/maverick-mavlink-router@.service
[Unit]
Description=Maverick mavlink-router Mavlink Proxy - %I
After=network-online.target
Requires=network-online.target

[Service]
TimeoutStartSec=0
EnvironmentFile=/srv/maverick/data/config/mavlink/mavlink-router-%I.service.conf
ExecStart=/srv/maverick/software/maverick/bin/mavlink-router.sh %I
User=mav

[Install]
WantedBy=multi-user.target