# manifests/maverick-modules/maverick_network/files/maverick-wifibc_tx.service
[Unit]
Description=Maverick WifiBC TX Daemon
After=network.target
Requires=network.target

[Service]
TimeoutStartSec=0
ExecStart=/srv/maverick/software/maverick/bin/wifibc.sh tx
User=mav

[Install]
WantedBy=multi-user.target

