# manifests/maverick-modules/maverick_ros/files/maverick-mavros@.service
[Unit]
Description=Maverick Mavros - %I
After=network-online.target
Requires=network-online.target

[Service]
TimeoutStartSec=0
ExecStart=/srv/maverick/software/maverick/bin/mavros.sh %I
User=mav

[Install]
WantedBy=multi-user.target
