# manifests/maverick-modules/maverick_ros/files/maverick-rosmaster@.service
[Unit]
Description=Maverick ROS Master - %I
After=maverick-motd.service
StartLimitBurst=5
After=multi-user.target

[Service]
TimeoutStartSec=0
ExecStart=/srv/maverick/software/maverick/bin/rosmaster.sh %I
User=mav
Restart=on-failure
RestartSec=30

[Install]
WantedBy=multi-user.target