# manifests/maverick-modules/maverick_ros/files/maverick-rosmaster@.service
[Unit]
Description=Maverick ROS Master - %I
After=maverick-mavproxy@%I.service maverick-mavlink-router@%I.service maverick-cmavnode@%I.service
Requires=network-online.target

[Service]
TimeoutStartSec=0
ExecStart=/srv/maverick/software/maverick/bin/rosmaster.sh %I
User=mav

[Install]
WantedBy=multi-user.target