# manifests/maverick-modules/maverick_dev/templates/dev-sitl.sim_vehicle.service.erb
[Unit]
Description=Maverick Dev ArduPilot SITL
After=network-online.target
Requires=network-online.target

[Service]
TimeoutStartSec=0
#ExecStart=/bin/bash -c 'cd /srv/maverick/.virtualenvs/sitl; source bin/activate; /srv/maverick/software/maverick/bin/dev-sitl.sh <%= @ardupilot_type %> >/srv/maverick/var/log/dev-sitl/dev-sitl.startup.log 2>&1'
ExecStart=/bin/bash -c 'cd /srv/maverick/.virtualenvs/sitl; source bin/activate; /usr/bin/screen -L -c /srv/maverick/data/config/dev-sitl.screen.conf -S dev-sitl -d -m /srv/maverick/code/ardupilot/Tools/autotest/sim_vehicle.py -N -v <%= @ardupilot_type %>'
User=mav
Type=forking

[Install]
WantedBy=multi-user.target

