# manifests/maverick-modules/maverick_dev/templates/maverick-sitl.service.erb
[Unit]
Description=Maverick ArduPilot SITL
After=network-online.target
Requires=network-online.target

[Service]
TimeoutStartSec=0
EnvironmentFile=/srv/maverick/data/config/dev/sitl.conf
ExecStart=/srv/maverick/software/maverick/bin/sitl.sh
User=mav

[Install]
WantedBy=multi-user.target

