# manifests/maverick-modules/maverick_dev/templates/mavproxy-sitl.service.erb
[Unit]
Description=Maverick Dev SITL Mavproxy
After=network-online.target maverick-dev-sitl.service
Requires=network-online.target maverick-dev-sitl.service

[Service]
TimeoutStartSec=0
Type=forking
ExecStart=/srv/maverick/software/maverick/bin/mavproxy-sitl.sh
User=mav

[Install]
WantedBy=multi-user.target

