# manifests/maverick-modules/maverick_dev/templates/sitl.service.erb
[Unit]
Description=Maverick Ardupilot SITL
After=network-online.target
Requires=network-online.target

[Service]
TimeoutStartSec=0
ExecStart=/bin/bash -c 'cd /srv/maverick/.virtualenvs/sitl; source bin/activate; ./bin/python ./bin/sitl /srv/maverick/code/ardupilot/<%= @sitl_fw_run %>/<%= @sitl_fw_run %>.elf --instance 1'
User=mav

[Install]
WantedBy=multi-user.target

