# manifests/maverick-modules/maverick_web/templates/cloud9.service.erb
[Unit]
Description=Cloud9
After=network-online.target
Requires=network-online.target

[Service]
TimeoutStartSec=0
ExecStart=/usr/bin/node /srv/maverick/software/cloud9/server.js -p <%= @webport %> --collab -l 0.0.0.0 -w <%= @basepath %> --packed --cache -a mav:<%= @cloud9_password %>
User=mav

[Install]
WantedBy=multi-user.target

