# manifests/maverick-modules/maverick_web/templates/codeserver.service.erb
[Unit]
Description=CodeServer
After=network-online.target
Requires=network-online.target

[Service]
TimeoutStartSec=0
WorkingDirectory=/srv/maverick
ExecStart=/srv/maverick/software/codeserver/packages/server/code-server --user-data-dir /srv/maverick/data/web/codeserver --port <%= @webport %> --allow-http --disable-telemetry --password <%= @password %> /srv/maverick
User=mav

[Install]
WantedBy=multi-user.target
