# manifests/maverick-modules/maverick_web/templates/codeserver.service.erb
[Unit]
Description=CodeServer
After=network-online.target
Requires=network-online.target

[Service]
TimeoutStartSec=0
WorkingDirectory=/srv/maverick
ExecStart=/srv/maverick/software/codeserver/codeserver --user-data-dir /srv/maverick/data/web/codeserver --extra-builtin-extensions-dir /srv/maverick/software/codeserver/extensions --extra-extensions-dir /srv/maverick/data/web/codeserver/extensions --port <%= @webport %> --disable-telemetry --disable-open --auth password /srv/maverick
User=mav

[Install]
WantedBy=multi-user.target
