# manifests/maverick-modules/maverick-dronekit/templates/dev-sitl.service.erb
[Unit]
Description=Maverick Dev Dronekit SITL
After=network-online.target
Requires=network-online.target

[Service]
TimeoutStartSec=0
ExecStart=/bin/bash -c 'cd /srv/maverick/.virtualenvs/dronekit-dev; source bin/activate; ./bin/python ./bin/dronekit-sitl /srv/maverick/code/dronekit-dev/ardupilot/<%= @sitl_fw_run %>/<%= @sitl_fw_run %>.elf --instance 1'
User=mav

[Install]
WantedBy=multi-user.target

