# manifests/maverick-modules/maverick-dronekit/templates/dev-mavproxy.service.erb
[Unit]
Description=Maverick Dev SITL Mavproxy
After=network-online.target
Requires=network-online.target dev-sitl.service

[Service]
TimeoutStartSec=0
ExecStart=/bin/bash -c 'cd /srv/maverick/.virtualenvs/dronekit-dev; source bin/activate; ./bin/mavproxy.py --master tcp:127.0.0.1:5770 --sitl 127.0.0.1:5501 --out=udpin:0.0.0.0:14560 --out=udpin:0.0.0.0:14561 --out=udp:0.0.0.0:14562 --out=udp:0.0.0.0:14563 --daemon'
User=mav

[Install]
WantedBy=multi-user.target

