# manifests/maverick-modules/maverick-dronekit/templates/fc-mavproxy.service.erb
[Unit]
Description=Maverick Flight Controller Mavproxy
After=network-online.target
Requires=network-online.target

[Service]
TimeoutStartSec=0
ExecStart=/bin/bash -c 'cd /srv/maverick/.virtualenvs/dronekit-fc; source bin/activate; ./bin/mavproxy.py <% if @fc_mavproxy_master %>--master <%= @fc_mavproxy_master %><% end %><% if @fc_mavproxy_baud %>--baud <%= @fc_mavproxy_baud %><% end %> --out=udpin:0.0.0.0:14550 --out=udpin:0.0.0.0:14551 --out=udp:127.0.0.1:14552 --out=udp:127.0.0.1:14553 --daemon'
User=mav

[Install]
WantedBy=multi-user.target

