# manifests/maverick-modules/maverick-network/templates/monitor-interface@.service.erb
[Unit]
Description=Monitor/injection mode interface setup for %I
After=network.target

[Service]
ExecStart=/srv/maverick/software/maverick/bin/monitor-interface.sh start %I
ExecStop=/srv/maverick/software/maverick/bin/monitor-interface.sh stop %I
User=root
RemainAfterExit=true

[Install]
WantedBy=multi-user.target

