# manifests/maverick-modules/base/files/maverick-motd.service
[Unit]
Description=Maverick MOTD Creator
After=network-online.target
Requires=network-online.target

[Service]
Type=oneshot
RemainAfterExit=True
TimeoutStartSec=0
ExecStart=/bin/sh -c '/srv/maverick/software/maverick/bin/maverick-info --motd >/etc/motd'
User=root

[Install]
WantedBy=multi-user.target

