# manifests/maverick-modules/maverick_vision/templates/maverick-visiond.service.erb
[Unit]
Description=Maverick Vision Daemon
After=network-online.target
Requires=network-online.target

[Service]
TimeoutStartSec=0
ExecStart=/srv/maverick/software/maverick/bin/maverick_visiond start
ExecStop=/srv/maverick/software/maverick/bin/maverick_visiond stop
User=mav

[Install]
WantedBy=multi-user.target

