# manifests/maverick-modules/maverick_cloud9/templates/cloud9.service.erb
[Unit]
Description=Cloud9
After=network-online.target
Requires=network-online.target

[Service]
TimeoutStartSec=0
ExecStart=/usr/bin/nodejs /srv/maverick/software/cloud9/server.js -p <%= @webport %> -l 0.0.0.0 -w <%= @basepath %> --packed --cache -a mav:wingman
User=mav

[Install]
WantedBy=multi-user.target

